# supervisor.conf
[unix_http_server]
file=/var/run/supervisor.sock              ; Unix套接字文件路径
chmod=0777                                 ; 套接字权限

[supervisord]
nodaemon=true                              ; 非守护模式，适合容器运行
pidfile=/tmp/supervisord.pid               ; PID文件路径

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; 使用Unix套接字与服务通信

[program:xiaomusic]
command=/app/.venv/bin/python3 /app/xiaomusic.py
directory=/app
autostart=true
autorestart=true
startretries=0
stderr_logfile=/app/xiaomusic.log.txt
stdout_logfile=/app/xiaomusic.log.txt

